# qwen3-server.service
[Unit]
Description=Qwen3-Coder-Next llama-server (dual RTX 3090)
After=network-online.target

[Service]
Type=simple
User=fabio
WorkingDirectory=/home/fabio/unsloth
Environment="LD_LIBRARY_PATH=/usr/local/cuda/lib64:/usr/lib/x86_64-linux-gnu"
ExecStart=/home/fabio/unsloth/llama.cpp/build/bin/llama-server \
  --model /home/fabio/unsloth/models/Qwen3-Coder-Next-IQ4_XS.gguf \
  --alias "unsloth/Qwen3-Coder-Next" \
  --host 0.0.0.0 \
  --port 8001 \
  --n-gpu-layers 999 \
  --split-mode layer \
  --tensor-split 1,1 \
  --ctx-size 65536 \
  --cache-type-k q8_0 \
  --cache-type-v q8_0 \
  --threads 16 \
  --batch-size 2048 \
  --temp 1.0 \
  --top-p 0.95 \
  --top-k 40 \
  --min-p 0.01 \
  --jinja \
  --seed 3407
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
