# worker-agent.service
[Unit]
Description=Worker Agent - Qwen3.5-27B-AWQ SGLang (4x RTX 3060)
After=network-online.target

[Service]
Type=simple
User=Mose
WorkingDirectory=/home/Mose/mose-agent
Environment="LD_LIBRARY_PATH=/usr/local/cuda/lib64:/usr/lib/x86_64-linux-gnu"
ExecStart=/home/Mose/mose-agent/.venv/bin/python -m sglang.launch_server \
  --model-path QuantTrio/Qwen3.5-27B-AWQ \
  --port 8001 \
  --tp-size 4 \
  --mem-fraction-static 0.85 \
  --context-length 98304 \
  --reasoning-parser qwen3 \
  --tool-call-parser qwen3_coder \
  --host 0.0.0.0
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
